# demo-reset.service
[Unit]
Description=Android demo session reset API
After=network.target docker.service

[Service]
Type=simple
ExecStart=/usr/bin/python3 /opt/android-demo/demo-reset-server.py
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
